# systemd/refract.target
[Unit]
Description=Refract target
After=systemd-modules-load.service

[Install]
WantedBy=multi-user.target
